# src/roles/host_bootstrap/files/disable-eth0-eee.service
[Unit]
Description=Disable Energy Efficient Ethernet (EEE) on eth0
After=network.target
Wants=network.target

[Service]
Type=oneshot
# Ignore exit code 76 (EOPNOTSUPP): Pi 5's NIC driver does not expose EEE
# controls via ethtool netlink; on Pi 4 the command succeeds and disables EEE.
ExecStart=/bin/sh -c '/usr/sbin/ethtool --set-eee eth0 eee off || true'
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=multi-user.target
